# armbian-build/userpatches/overlay/etc/systemd/system/easytier.service
[Unit]
Description=EasyTier Mesh Network
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/start-easytier.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
