# armbian-build/userpatches/overlay/etc/systemd/system/hive-clash-sync.timer
[Unit]
Description=Run hive-clash-sync periodically

[Timer]
# 启动后 2min 首次运行（给 mihomo/网络留就绪时间），之后每 60s 一次
# 带 15s 抖动，避免所有节点同时打 registry
OnBootSec=2min
OnUnitActiveSec=60s
RandomizedDelaySec=15s
AccuracySec=5s

[Install]
WantedBy=timers.target
